# linux/files/fortress-update.service
[Unit]
Description=Fortress IaC - Adult Shield & SafeSearch Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fortress-update
