# linux/files/fortress-update.timer
[Unit]
Description=Fortress IaC - Shield Update Trigger

[Timer]
# Runs 5 minutes after boot to ensure network has stabilized
OnBootSec=5min
# Runs once a week (you can change to 'daily' if preferred)
OnCalendar=weekly
# If the machine was off at the scheduled time, runs on next boot
Persistent=true

[Install]
WantedBy=timers.target
